# data/useless-logind.service.in
[Unit]
Description=Useless Login Service

[Service]
ExecStart=@systemdsystemunitdir@/useless-logind
Restart=always
RestartSec=0
